# .sisyphus/notepads/fzs-modernization/supervisor-config.ini
; Supervisor konfiguracija za FZS Laravel queue workers
; Putanja: /etc/supervisor/conf.d/fzs-laravel-worker.conf
;
; Instalacija i podešavanje:
;   sudo apt-get install supervisor
;   sudo cp fzs-laravel-worker.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start fzs-laravel-worker:*
;
; Praćenje statusa:
;   sudo supervisorctl status
;   sudo supervisorctl tail -f fzs-laravel-worker:fzs-laravel-worker_00

[program:fzs-laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/dragomix/fzs-laravel/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=dragomix
numprocs=2
redirect_stderr=true
stdout_logfile=/home/dragomix/fzs-laravel/storage/logs/worker.log
stopwaitsecs=3600

; Napomena: stopwaitsecs=3600 omogućava MassEnrollmentJob (timeout=600) da se dovrši
; pre nego što Supervisor prinudno zaustavi worker.
; numprocs=2 daje dva paralelna workera za veće opterećenje.
; Za produkciju razmotriti povećanje numprocs na 4-8 u zavisnosti od servera.
